# scripts/openrc/sanbox.sh
#!/sbin/openrc-run

command_user="sanbox:docker"
description="dockers rootless cleanup service"
